# deploy/pi-sense-agent.service
[Unit]
Description=SenseHAT MQTT Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rpicd/SenseHat_ConditionMonitor
Environment=MQTT_HOST=192.168.1.10
Environment=MQTT_PORT=1883
Environment=MQTT_TLS=0
ExecStart=/home/rpicd/SenseHat_ConditionMonitor/edge/.venv/bin/python /home/rpicd/SenseHat_ConditionMonitor/edge/agent.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
